# provision/etc/supervisor/conf.d/myapp.conf
[program:myapp]
directory=/srv/myapp/current/myapp
command=/srv/myapp/shared/env/bin/gunicorn_django -c /etc/gunicorn/myapp.conf.py
user=myapp
autostart=true
autorestart=true
redirect_stderr=True
